# data/signaling/janus.service
[Unit]
Description=Janus WebRTC Server
Documentation=https://janus.conf.meetecho.com/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Janus 1.2+: -F is --configs-folder (path required); it is not "foreground".
# Omit --daemon/-b so Janus stays in the foreground for Type=simple.
ExecStart=/usr/bin/janus -F /etc/janus
Restart=on-failure
RestartSec=3s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
